# cron/etc/cron.daily/renew-cert
#!/bin/bash
DIRNAME="$(readlink -f $0)";
cd "$(dirname "${DIRNAME}")"/../../../
./certbot-renew.sh;
docker-compose exec nginx nginx -s reload;

